# supervisord.conf
[supervisord]
nodaemon=true

[program:icecast2]
command=/usr/bin/icecast2 -c /etc/icecast2/icecast.xml
user=icecast2
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
